# 07-headless-jamulus/files/update.service
[Unit]
Description=login update jamulus

[Service]
ExecStart=/home/pi/autoupdate.sh
Type=oneshot

[Install]
WantedBy=default.target

